# os/rootfs-overlay/etc/systemd/system/rocola.service
# Lanza la sesión de la rocola: Xorg mínimo + app SDL2 a pantalla completa en tty1.
#
# Patrón "reemplazo de getty en tty1" (kiosko): en vez de competir con getty@tty1
# vía Conflicts (lo que hacía que systemd descartara silenciosamente uno de los
# dos), este servicio toma tty1 directamente. getty@tty1 se deshabilita por
# separado enmascarándolo (ver getty@tty1 mask en el overlay) para que nada más
# reclame la VT1.
[Unit]
Description=Rocola (UI SDL2 + Xorg en tty1)
# Arrancar después de que el sistema esté listo, mpd levantado y el video configurado.
After=systemd-user-sessions.service plymouth-quit-wait.service mpd.service rocola-display-setup.service
Wants=mpd.service rocola-display-setup.service
# Rate-limit: reintentar ante fallos transitorios, pero RENDIRSE si el fallo es
# permanente (p. ej. Xorg no levanta por falta de driver/firmware). Sin esto, con
# Restart=always + StartLimitIntervalSec=0 el servicio flapea infinito cada 2 s,
# roba el VT/pantalla y deja la máquina en un login-loop inusable. Tras 5 fallos
# en 120 s systemd lo marca 'failed' y NO reintenta más: tty1 cae a un login
# normal (autovt) donde se puede leer 'journalctl -u rocola.service' y depurar.
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
User=rocola
PAMName=login
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
StandardInput=tty
StandardOutput=journal
StandardError=journal
# xinit arranca Xorg en vt1 y, como cliente único, la sesión de la rocola.
ExecStart=/usr/bin/xinit /usr/local/bin/rocola-session -- :0 vt1 -keeptty -nolisten tcp
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
