# os/rootfs-overlay/etc/systemd/system/getty@tty1.service.d/autologin.conf
# Autologin del usuario `rocola` en tty1 (sin contraseña).
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin rocola --noclear %I $TERM
