# os/rootfs-overlay/etc/systemd/system/getty@tty1.service.d/autologin.conf
# Autologin del usuario 'rocola' en tty1 (patrón kiosko getty-autologin + startx).
#
# Por qué getty-autologin y NO un servicio propio que reemplace getty:
# la sesión que crea el login del autologin queda registrada en logind como la
# sesión ACTIVA de seat0/vt1. Así, cuando la UI arranca Xorg (vía startx en
# ~/.bash_profile), logind le entrega el fd del DRM en vez de un "paused fd". El
# patrón anterior (servicio tomando tty1 con xinit) chocaba con logind: Xorg
# recibía "paused fd for drm node" y/o xinit moría por SIGHUP. Ver docs/03.
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin rocola --noclear %I linux
